# deploy/ansible/roles/services/templates/tessel-io.conf
description "The Tessel website"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}
env PORT={{app_port}}

pre-start script
  echo "[`date +%Y%m%d%H%M%S`] starting"
end script

script
  /usr/bin/node {{site_path}}/index
end script

pre-stop script
  echo "[`date +%Y%m%d%H%M%S`] stopping "
end script
